# stt/song_video/holo.song.video_holo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>holo.song.video_holo</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Library/Apple/usr/bin</string>
    </dict>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/lua</string>
      <string>/Users/kamui/doc/hoby/youtube/vtuber/holo/prj/stt/song_video/song_video_holo.lua</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kamui/doc/hoby/youtube/vtuber/holo/prj/stt/song_video</string>
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Weekday</key>
        <integer>4</integer>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>45</integer>
      </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/kamui/doc/hoby/youtube/vtuber/holo/prj/stt/song_video/log/holo.std-out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kamui/doc/hoby/youtube/vtuber/holo/prj/stt/song_video/log/holo.std-err.log</string>
  </dict>
</plist>

